# 10-generic-radar.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="04e2", ATTRS{idProduct}=="1412", RUN+="/usr/local/bin/radar_load.sh"
